# ctfs/TetCTF/2023/pwn/MailService/shares/mailserver.service
[Unit]
Description=Mailserver Service

[Service]
User=mailserver
WorkingDirectory=/home/mailserver/
ExecStart=/home/mailserver/mailserver
Restart=always

[Install]
WantedBy=multi-user.target