# ctfs/BraekerCTF/2024/web/Stuffy/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx

[program:python]
command=python app.py
