# ctfs/Hack.lu/2023/web/Mastodont/mastodon/hacklu-init/etc/s6-overlay/s6-rc.d/setup-challenge/up
/etc/s6-overlay/s6-rc.d/setup-challenge/run