# HITCON/2022/pwn/EaaS/clear.cron
*/1 * * * * /bin/bash -c "rm -f /tmp/eaas*"
