# ctfs/SekaiCTF/2023/web/Scanner_Service/config/supervisord.conf
[supervisord]
user=root
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[program:app]
directory=/app
user=www
command=bundle exec rackup -o0.0.0.0 -p1337
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0