# ctfs/TCP1P-CTF/2023/web/PDFIFY/src/infra/docker/supervisord/apache.conf
[program:apache]
command = apachectl -D "FOREGROUND" -k start
stdout_logfile = /dev/stdout
stderr_logfile = /dev/stderr
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0
autostart = true
user = root
autorestart = true
