# ctfs/corCTF/2023/web/msfrognymize/supervisord.conf

[supervisord]
nodaemon=true

[program:redis]
command=redis-server --bind 0.0.0.0
autostart=true
autorestart=true

[program:gunicorn]
command=gunicorn --bind 0.0.0.0:4444 --workers 4 app:app
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/gunicorn.err.log
stdout_logfile=/var/log/gunicorn.out.log

[program:celery]
command=celery -A tasks.celery_app worker --loglevel=info
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/celery.err.log
stdout_logfile=/var/log/celery.out.log

[program:celery-beat]
command=celery -A tasks.celery_app beat --scheduler=redbeat.RedBeatScheduler --loglevel=info
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/celery-beat.err.log
stdout_logfile=/var/log/celery-beat.out.log
