# CSAW/2021/Quals/web/Gatekeeping/supervisord.conf
[supervisord]
nodaemon=true

[program:gunicorn3]
user=app
command=gunicorn -w 10 -b unix:/tmp/gunicorn.sock --user app server:app
directory=/server
autorestart=true

[program:nginx]
user=root
command=nginx -g 'daemon off;'
directory=/
autorestart=true